# contrib/xz-mirror.service
[Unit]
Description=XZ Utils Source Code Mirroring Script

[Service]
Type=oneshot
ExecStart=/home/simnalamburt/services/xz-mirror-script/src/do
